# references/mORMot/SQLite3/Samples/36 - Simple REST Benchmark/rest_benchmark.socket
[Socket]
ListenStream=127.0.0.1:8889
NoDelay=true
KeepAlive=true
 
[Install]
WantedBy=sockets.target